# data/component-tags.json
{
  "description": {
    "all": [
      "sql",
      "xml",
      "web",
      "security",
      "database",
      "json",
      "yaml",
      "validation",
      "sanitization",
      "cloud",
      "iam",
      "auth",
      "middleware",
      "serialization",
      "event",
      "stream",
      "rpc",
      "socket",
      "proto",
      "resource",
      "sensitive",
      "template",
      "log",
      "logging",
      "service",
      "api",
      "slf4j",
      "parse",
      "emit",
      "jdbc",
      "connect",
      "pool",
      "beans",
      "transaction",
      "mysql",
      "postgres",
      "oracle",
      "mongo",
      "redis",
      "splunk",
      "stripe",
      "payment",
      "finance",
      "currency",
      "coin",
      "monero",
      "ssl",
      "traffic",
      "mvc",
      "html",
      "escape",
      "unescape",
      "rest",
      "tomcat",
      "hibernate",
      "orm",
      "aop",
      "jwt",
      "saml",
      "token",
      "tls",
      "codec",
      "cron",
      "crypto",
      "jce",
      "certificate",
      "developer",
      "tools",
      "autoconfigure",
      "test",
      "jsonpath",
      "bytecode",
      "mock",
      "injection",
      "comparators",
      "transform",
      "encode",
      "decode",
      "ldap",
      "owasp",
      "fileupload",
      "beanshell",
      "spel",
      "mail",
      "apacheds",
      "jndi",
      "ldif",
      "jdbm",
      "kerberos",
      "oidc",
      "oauth2",
      "cli",
      "binary",
      "ml",
      "ai",
      "azure",
      "gcp",
      "terraform",
      "redis",
      "valkey",
      "lint",
      "bundle",
      "object-persistence",
      "text-to-image",
      "translat",
      "object-detect",
      "mvc",
      "framework",
      "graph",
      "templates",
      "fastjson",
      "simd",
      "event-driven",
      "productivity",
      "typesafe",
      "projections",
      "performance",
      "plugins",
      "non-block"
    ]
  },
  "properties": {
    "all": [
      "sql",
      "http",
      "xml",
      "cloud",
      "middleware",
      "framework",
      "bluetooth",
      "wifi",
      "wireless",
      "driver",
      "graphics",
      "firmware",
      "gyroscope",
      "accelerometer",
      "mobile",
      "network",
      "battery",
      "matrix",
      "thunderbolt",
      "crypto",
      "algorithm",
      "encrypt",
      "decrypt",
      "registry",
      "payment",
      "stripe",
      "apple-pay",
      "icloud"
    ],
    "obom": [
      "windows_drivers",
      "windows_patches",
      "windows_programs",
      "processor",
      "services_snapshot",
      "apt_sources",
      "behavioral_reverse_shell",
      "certificates",
      "chrome_extensions",
      "crontab_snapshot",
      "deb_packages",
      "docker_container_ports",
      "docker_containers",
      "docker_networks",
      "docker_volumes",
      "etc_hosts",
      "firefox_addons",
      "vscode_extensions",
      "homebrew_packages",
      "installed_applications",
      "interface_addresses",
      "kernel_info",
      "kernel_integrity",
      "kernel_modules",
      "ld_preload",
      "listening_ports",
      "os_version",
      "pipes",
      "pipes_snapshot",
      "portage_packages",
      "process_events",
      "processes",
      "python_packages",
      "rpm_packages",
      "scheduled_tasks",
      "services_snapshot",
      "startup_items",
      "system_info_snapshot",
      "windows_drivers",
      "windows_patches",
      "windows_programs",
      "windows_shared_resources",
      "yum_sources",
      "appcompat_shims",
      "browser_plugins",
      "certificates",
      "chocolatey_packages",
      "chrome_extensions",
      "etc_hosts",
      "firefox_addons",
      "ie_extensions",
      "kernel_info",
      "npm_packages",
      "opera_extensions",
      "pipes_snapshot",
      "process_open_sockets",
      "safari_extensions",
      "scheduled_tasks",
      "services_snapshot",
      "startup_items",
      "routes",
      "system_info_snapshot",
      "win_version",
      "windows_firewall_rules",
      "windows_optional_features",
      "windows_programs",
      "windows_shared_resources",
      "windows_update_history",
      "wmi_cli_event_consumers",
      "wmi_cli_event_consumers_snapshot",
      "wmi_event_filters",
      "wmi_filter_consumer_binding"
    ]
  },
  "name": {
    "sbom": [
      {
        "test": [
          "(junit|xmlunit|testng|mocha|jest|test4j|xunit|coverlet|Test\\.Sdk)",
          "^(chai)$"
        ]
      },
      {
        "security": [
          "(boringssl|openssl|libressl|libssl|gnutls|jose|keyutils|Azure\\.Security|System\\.Security)"
        ]
      },
      { "native": ["(ffi|native)"] },
      { "parse": ["(parser)"] },
      { "transform": ["(transformer)"] },
      { "telemetry": ["(OpenTelemetry)"] },
      { "logging": ["(Microsoft\\.Extensions\\.Logging|Log4net)"] },
      { "ml": ["^(llama|langchain|openai)", "(generativeai)"] },
      { "devel": ["^(types-|typing-|virtualenv|ruff|poetry)"] }
    ],
    "obom": [
      {
        "devel": [
          "-(dev|devel|headers|sdk|libs|extension|headers+x86|headers+x64|headers+arm64)$",
          "^(git)[-]?",
          "^(sdk|windows+sdk)"
        ]
      },
      {
        "bin": [
          "(-bin|redistributable|clickonce|bootstrappermsi|bootstrappermsires|clickoncesigntoolmsi|codecoveragemsi|msires|sharedmsi|x64msi|arm64msi|sharedmsi|x64vmsi|filehandler_amd64|filehandler_x86|protocolhandlermsi|interopmsi|interopx64msi|shellmsires|shellx64msi)$"
        ]
      },
      { "kernel": ["^(linux|kernel|os-image)"] },
      {
        "security": [
          "(selinux|apparmor|security|boringssl|openssl|libressl|gnutls|jose|keyutils|passwd|libssl|libaudit|gcrypt|libpam|libseccomp)"
        ]
      },
      {
        "container": [
          "(container|podman|docker|runc|nerdctl|crun|libvirt|qemu)"
        ]
      },
      {
        "build": [
          "(cpp|fortran|gcc|make|meson|bazel|maven|gradle|sbt|ant|gdb|boost|compiler|kotlin|cargo|rustc|llvm|libstdc)"
        ]
      },
      {
        "network": [
          "(tailscale|wireguard|openvpn|dns|cockpit|cups|dhcp|network|iproute|iptables|mosh|netavark|openssh|rsync|tcpdump|libssh)"
        ]
      },
      { "webserver": ["(httpd|http2|tomcat|jboss)"] },
      {
        "crypto": [
          "(crypt|gpg|keys|certificates|gnupg|certifi|pubkey|keyutils|nss|keyring)"
        ]
      },
      { "repository": ["(-repos|-release|ostree|appstream)"] },
      { "shell": ["(bash|zsh|csh|fish|binsh|dash|oilsh)"] },
      { "bluetooth": ["(bluez|bluetooth)"] },
      { "sound": ["(alsa|pulseaudio|wireplumber|flac|codecs|ldac|sound)"] },
      {
        "compression": [
          "(brotli|xz-utils|zstd|lz4|zlib|bz2|lzma5|bzip2|libarchive)",
          "(tar|zip|webp|zchunk)$"
        ]
      },
      {
        "runtime": [
          "(perl|lua|php|python|ruby|dotnet|java|swift|runtime|glibc|libc6|musl|wasm|\\.net|asp\\.net|node.js|node|groovy)"
        ]
      },
      { "editor": ["(vim|emacs|nano|hexedit)"] },
      { "xml": ["(xml|expat)"] },
      { "boot": ["(grub|systemd-boot|syslinux|init-system|sysvinit)"] },
      {
        "gui": [
          "(wayland|xorg|X11|mesa|vulkan|tk|wkhtmltox|electron|Xrender|glib2)"
        ]
      },
      {
        "package": [
          "(rpm|dnf|yum|apt|zypper|apk|conda|dpkg|dnf5)$",
          "^(conda_package_|conda-package-|libapt|dnf5-|libdnf5)"
        ]
      },
      {
        "browser": [
          "^(edge)",
          "(firefox|chrome|opera|brave|mullvad|tor|chromium)",
          "(microsoft+edge|microsoft+edge+webview2|microsoft+html)"
        ]
      },
      {
        "chat": ["(webex|teams|slack|discord|vesktop|matrix|signal|whatsapp)"]
      },
      { "logging": ["(log4j|logging|slf4j)"] },
      { "root": ["^(sudo|systemd|pam|shadow)$"] }
    ]
  }
}
